# sv/fluxbox/run
#!/bin/bash

sv start xvfb || exit 1

. /etc/settings.env

exec /usr/bin/fluxbox -display "$DISPLAY" -screen 0
