# rootfs/etc/sv/xvfb/run
#!/bin/bash

# no dependencies, can start immediately

. /etc/settings.env

exec Xvfb "$DISPLAY" -listen tcp -screen 0 "$SCREEN"
