# rootfs/etc/sv/fluxbox/run
#!/bin/bash

# specify dependency is xvfb
# must run X server before window manager
sv start xvfb || exit 1

. /etc/settings.env

exec /usr/bin/fluxbox -display "$DISPLAY" -screen 0
