# rootfs/etc/sv/x11vnc/run
#!/bin/bash

sv start xvfb || exit 1

. /etc/settings.env

exec x11vnc -rfbport "$RFBPORT" -display "$DISPLAY" -N -forever
