# sv/x11vnc/run
#!/bin/bash

sv start x11vnc || exit 1

. /etc/settings.env

exec x11vnc -rfbport "$RFBPORT" -display "$DISPLAY" -N -forever
